# systemd/akai-magicq-bridge.service
[Unit]
Description=AKAI APC Mini MK2 MagicQ Bridge
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/bridge/akai-magicq-bridge
ExecStart=/usr/bin/node server/index.js
Restart=always
User=akai
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
